# backend/deploy/supervisor/finengine.conf
[group:finengine]
programs=finengine-api,finengine-worker,finengine-flower

[program:finengine-api]
directory=/var/www/FinEngine/backend
command=/bin/bash /var/www/FinEngine/backend/deploy/start_api.sh
user=www
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=60
stopasgroup=true
killasgroup=true
stopsignal=TERM
stdout_logfile=/var/www/FinEngine/backend/logs/supervisor/api.stdout.log
stderr_logfile=/var/www/FinEngine/backend/logs/supervisor/api.stderr.log
stdout_logfile_maxbytes=20MB
stderr_logfile_maxbytes=20MB
stdout_logfile_backups=10
stderr_logfile_backups=10
environment=PYTHONUNBUFFERED="1",APP_HOST="127.0.0.1",APP_PORT="8000",APP_WORKERS="2",APP_FORWARDED_ALLOW_IPS="*"

[program:finengine-worker]
process_name=%(program_name)s_%(process_num)02d
numprocs=2
directory=/var/www/FinEngine/backend
command=/bin/bash /var/www/FinEngine/backend/deploy/start_worker.sh
user=www
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=120
stopasgroup=true
killasgroup=true
stopsignal=TERM
stdout_logfile=/var/www/FinEngine/backend/logs/supervisor/worker_%(process_num)02d.stdout.log
stderr_logfile=/var/www/FinEngine/backend/logs/supervisor/worker_%(process_num)02d.stderr.log
stdout_logfile_maxbytes=20MB
stderr_logfile_maxbytes=20MB
stdout_logfile_backups=10
stderr_logfile_backups=10
environment=PYTHONUNBUFFERED="1",CELERY_POOL="solo",CELERY_HOSTNAME="finengine-%(process_num)02d@%%h"

[program:finengine-flower]
directory=/var/www/FinEngine/backend
command=/bin/bash /var/www/FinEngine/backend/deploy/start_flower.sh
user=www
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=30
stopasgroup=true
killasgroup=true
stopsignal=TERM
stdout_logfile=/var/www/FinEngine/backend/logs/supervisor/flower.stdout.log
stderr_logfile=/var/www/FinEngine/backend/logs/supervisor/flower.stderr.log
stdout_logfile_maxbytes=20MB
stderr_logfile_maxbytes=20MB
stdout_logfile_backups=10
stderr_logfile_backups=10
environment=PYTHONUNBUFFERED="1",FLOWER_HOST="127.0.0.1",FLOWER_PORT="5555",FLOWER_BASIC_AUTH="",FLOWER_URL_PREFIX="flower"
