# backend/deploy/supervisor/finengine.conf
[group:finengine]
programs=finengine-api,finengine-worker

[program:finengine-api]
directory=/var/www/FinEngine/backend
command=/bin/bash /var/www/FinEngine/backend/deploy/start_api.sh
user=www
autostart=true
autorestart=true
startsecs=5
stopasgroup=true
killasgroup=true
stopsignal=TERM
stdout_logfile=/var/www/FinEngine/backend/logs/supervisor/api.stdout.log
stderr_logfile=/var/www/FinEngine/backend/logs/supervisor/api.stderr.log
stdout_logfile_maxbytes=20MB
stderr_logfile_maxbytes=20MB
stdout_logfile_backups=10
stderr_logfile_backups=10
environment=PYTHONUNBUFFERED="1"

[program:finengine-worker]
directory=/var/www/FinEngine/backend
command=/bin/bash /var/www/FinEngine/backend/deploy/start_worker.sh
user=www
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=60
stopasgroup=true
killasgroup=true
stopsignal=TERM
stdout_logfile=/var/www/FinEngine/backend/logs/supervisor/worker.stdout.log
stderr_logfile=/var/www/FinEngine/backend/logs/supervisor/worker.stderr.log
stdout_logfile_maxbytes=20MB
stderr_logfile_maxbytes=20MB
stdout_logfile_backups=10
stderr_logfile_backups=10
environment=PYTHONUNBUFFERED="1"
